# files/cron.d/puppetmaster
# we restart puppetmaster twice a day to avoid memory problems
58 11 * * * root /etc/init.d/puppetmaster restart > /dev/null

