# files/cron.d/puppetd.CentOS
0 * * * *   root    /bin/ps ax | /bin/grep -v grep | /bin/grep puppetd > /dev/null || /etc/init.d/puppet restart
