# files/cron.d/puppetmaster.CentOS
# we restart puppetmaster every 2 hours to avoid memory problems
0 */2 * * * root /etc/init.d/puppetmaster restart > /dev/null
