# files/cron.d/puppetmaster
# we restart puppetmaster every 4 hours to avoid memory problems
0 */6 * * * root /etc/init.d/puppetmaster restart > /dev/null

