# files/cron.d/puppetmaster
# we restart puppetmaster twice a day to avoid memory problems
59 11,23 * * * root /etc/init.d/puppetmaster restart > /dev/null

